# scripts/wazuh-auto-dls.service
[Unit]
Description=Auto-manage per-agent DLS roles/users for Wazuh
After=network.target wazuh-manager.service wazuh-indexer.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wazuh-auto-dls.sh
